# docker/orbit-emacs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_VNC_GEOMETRY)sx%(ENV_VNC_DEPTH)s -nolisten tcp
user=orbit
priority=10
autorestart=true
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb.err

[program:openbox]
command=/usr/bin/openbox
user=orbit
priority=20
autorestart=true
environment=DISPLAY=":1",HOME="/home/orbit"
stdout_logfile=/tmp/openbox.log
stderr_logfile=/tmp/openbox.err

[program:vnc]
command=/usr/local/bin/orbit-start-vnc
user=orbit
priority=30
autorestart=true
environment=DISPLAY=":1",HOME="/home/orbit"
stdout_logfile=/tmp/x11vnc.log
stderr_logfile=/tmp/x11vnc.err

[program:emacs]
command=/usr/local/bin/orbit-start-emacs
user=orbit
priority=40
autorestart=false
startsecs=3
environment=DISPLAY=":1",HOME="/home/orbit",ORBIT_WORKSPACE="/workspace"
stdout_logfile=/tmp/emacs.log
stderr_logfile=/tmp/emacs.err
