# docker/orbit-emacs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:vnc]
command=/usr/local/bin/orbit-start-vnc
user=orbit
priority=10
autorestart=true
environment=DISPLAY=":1",HOME="/home/orbit",ORBIT_HOME="/home/orbit"
stdout_logfile=/tmp/vnc.log
stderr_logfile=/tmp/vnc.err

[program:openbox]
command=/usr/local/bin/orbit-start-window-manager
user=orbit
priority=20
autorestart=true
environment=DISPLAY=":1",HOME="/home/orbit",ORBIT_HOME="/home/orbit"
stdout_logfile=/tmp/openbox.log
stderr_logfile=/tmp/openbox.err

[program:emacs]
command=/usr/local/bin/orbit-start-emacs
user=orbit
priority=30
autorestart=false
startsecs=3
environment=DISPLAY=":1",HOME="/home/orbit",ORBIT_HOME="/home/orbit",ORBIT_WORKSPACE="/workspace"
stdout_logfile=/tmp/emacs.log
stderr_logfile=/tmp/emacs.err
